# python/api.conf
[program:api]
directory=/root/
command=gunicorn api:api -b 0.0.0.0:8080
autostart=true
autorestart=true
stderr_logfile=/var/log/python-api.err.log
stdout_logfile=/var/log/python-api.out.log